# docker/supervisor/conf.d/queues.conf
# 消息队列消费者进程配置

[program:withdraw_queue]
command=php /var/www/html/think queue:listen --queue withdraw
process_name=%(program_name)s_%(process_num)02d
numprocs=3
autostart=true
autorestart=true
startsecs=1
startretries=3
stopwaitsecs=30
stopsignal=QUIT
stdout_logfile=/var/log/supervisor/withdraw_queue.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stderr_logfile=/var/log/supervisor/withdraw_queue_error.log
user=www-data
priority=100

[program:video_reward_queue]
command=php /var/www/html/think queue:listen --queue video_reward
process_name=%(program_name)s_%(process_num)02d
numprocs=5
autostart=true
autorestart=true
startsecs=1
startretries=3
stopwaitsecs=30
stopsignal=QUIT
stdout_logfile=/var/log/supervisor/video_reward_queue.log
stderr_logfile=/var/log/supervisor/video_reward_queue_error.log
user=www-data
priority=100

[program:redpacket_queue]
command=php /var/www/html/think queue:listen --queue redpacket
process_name=%(program_name)s_%(program_name)s_%(process_num)02d
numprocs=3
autostart=true
autorestart=true
startsecs=1
startretries=3
stopwaitsecs=30
stopsignal=QUIT
stdout_logfile=/var/log/supervisor/redpacket_queue.log
stderr_logfile=/var/log/supervisor/redpacket_queue_error.log
user=www-data
priority=100

[program:commission_queue]
command=php /var/www/html/think queue:listen --queue commission
process_name=%(program_name)s_%(process_num)02d
numprocs=2
autostart=true
autorestart=true
startsecs=1
startretries=3
stopwaitsecs=30
stopsignal=QUIT
stdout_logfile=/var/log/supervisor/commission_queue.log
stderr_logfile=/var/log/supervisor/commission_queue_error.log
user=www-data
priority=100

[program:notification_queue]
command=php /var/www/html/think queue:listen --queue notification
process_name=%(program_name)s_%(process_num)02d
numprocs=2
autostart=true
autorestart=true
startsecs=1
startretries=3
stopwaitsecs=30
stopsignal=QUIT
stdout_logfile=/var/log/supervisor/notification_queue.log
stderr_logfile=/var/log/supervisor/notification_queue_error.log
user=www-data
priority=100

[program:default_queue]
command=php /var/www/html/think queue:listen --queue default
process_name=%(program_name)s_%(process_num)02d
numprocs=3
autostart=true
autorestart=true
startsecs=1
startretries=3
stopwaitsecs=30
stopsignal=QUIT
stdout_logfile=/var/log/supervisor/default_queue.log
stderr_logfile=/var/log/supervisor/default_queue_error.log
user=www-data
priority=70
